# .config/autostart/firefox.desktop
[Desktop Entry]
Version=1.0
Name=Firefox
Exec=firefox
Terminal=false
Type=Application
Icon=firefox
StartupNotify=true
